# scripts/templates/etc/inittab
::sysinit:/etc/init.d/rcS
::shutdown:/bin/umount -a -r
ttyS0::respawn:/sbin/getty 115200 ttyS0